# src/crontab
@reboot python3 /home/farm/Documents/the-python-mega-course/src/webBlocker.py

